# scripts/com.serydoo.photomemo.autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.serydoo.photomemo.autosync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/rui/Desktop/PhotoMemo/scripts/auto_sync_to_github.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>WorkingDirectory</key>
  <string>/Users/rui/Desktop/PhotoMemo</string>

  <key>StandardOutPath</key>
  <string>/Users/rui/Desktop/PhotoMemo/.autosync/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/rui/Desktop/PhotoMemo/.autosync/launchd.err.log</string>
</dict>
</plist>
