# scripts/com.serydoo.photomemo.autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.serydoo.photomemo.autosync</string>

  <key>ProgramArguments</key>
  <array>
    <string>__TARGET_SCRIPT__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>WorkingDirectory</key>
  <string>__SUPPORT_DIR__/AutoSync</string>

  <key>StandardOutPath</key>
  <string>__SUPPORT_DIR__/AutoSync/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>__SUPPORT_DIR__/AutoSync/launchd.err.log</string>
</dict>
</plist>
